# systemd/honeytoken.service
[Unit]
Description=HoneyToken service

[service]
ExecStart=/home/regular/Documents/projects/canarytokens/eh.py
