# systemd/honeytoken.service
[Unit]
Description=HoneyToken service

[service]
ExecStart=/opt/honeytokens/honeytokens.py
